# 99-thinksmart-ha.rules
# udev rules for the ThinkSmart Hub 500 lnxlink custom modules.
# Grant the session user device access without root:
#   - USB 17ef:a017 (LED controller, vendor ep0 transfers)
#   - PIR IIO attributes (arming sampling_frequency)
# Install into /etc/udev/rules.d/ then: udevadm control --reload && udevadm trigger

# USB card / LED controller -- usbfs node access (no interface detach)
SUBSYSTEM=="usb", ATTR{idVendor}=="17ef", ATTR{idProduct}=="a017", MODE="0660", GROUP="plugdev", TAG+="uaccess"

# PIR sensor exposed via IIO (hid-sensor-prox) -- make sampling_frequency/raw
# group-writable (its name is "prox").
SUBSYSTEM=="iio", KERNEL=="iio:device*", ATTR{name}=="prox", MODE="0664", GROUP="plugdev", TAG+="uaccess", \
  RUN+="/bin/sh -c 'chgrp -R plugdev /sys%p/ 2>/dev/null; chmod -R g+rw /sys%p/in_proximity0_sampling_frequency /sys%p/in_proximity0_raw 2>/dev/null'"

# Screen: DDC access goes through /dev/i2c-5 (i2c group, handled by i2c-tools rules).
# Make sure the user is in the i2c group (see install.sh).
